# package/base-files/files/etc/crontabs/root
# 10 4 * * * reboot
*/10 * * * * /sbin/ntptohc.sh
