# feeds/luci/themes/luci-theme-entels/root/etc/crontabs/root
*/10 * * * * /sbin/ntptohc.sh
*/10 * * * * /sbin/remote_upgrade
