# package/base-files/files/etc/udev/rules.d/30-rtc.rules
KERNEL=="rtc0",RUN+="/usr/sbin/hwclock --localtime -s"