# package/base-files/files/etc/hotplug.d/usb/31-cp2114
if [ "$ACTION" == "add" -a "$PRODUCT" == "10c4/ea90/0" ]; then
    echo 507 > /sys/class/gpio/export
    echo 508 > /sys/class/gpio/export
    echo high > /sys/class/gpio/gpio508/direction
    echo low > /sys/class/gpio/gpio507/direction
    echo high > /sys/class/gpio/gpio507/direction
fi
