# package/base-files/files/etc/rc.local
/etc/initscripts/rtcinit &
/etc/initscripts/ethinit &
/etc/initscripts/usbhost &
/etc/initscripts/modeminit &
/etc/initscripts/pciinit &

watchdog -t 30 -T 60 /dev/watchdog0

exit 0
