# package/base-files/files/etc/rc.local
echo 105 > /sys/class/gpio/export
echo out > /sys/class/gpio/gpio105/direction
echo 1 > /sys/class/gpio/gpio105/value

echo 83 > /sys/class/gpio/export
echo out > /sys/class/gpio/gpio83/direction
echo 0 > /sys/class/gpio/gpio83/value

/etc/initscripts/usbhost
watchdog -t 30 -T 60 /dev/watchdog0
exit 0
