# package/base-files/files/etc/rc.local
/etc/initscripts/usbhost
watchdog -t 30 -T 60 /dev/watchdog0
exit 0
