# feeds/luci/themes/luci-theme-entels/root/etc/uci-defaults/luci-theme-entels
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Entels=/luci-static/entels
	set luci.main.mediaurlbase=/luci-static/entels
	commit luci
EOF
exit 0
